# container-files/etc/services.d/in.tftpd/run
#!/bin/execlineb -P

/usr/sbin/in.tftpd --port-range 4096:32767 -m /var/lib/tftpboot/pxelinux.cfg/remap --foreground --user tftpd --secure /data
